# infrastructure/retail-generator.service
[Unit]
Description=Retail Store Performance Data Generator
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=telco_user
Group=sftpusers
WorkingDirectory=/opt/telco-generator
Environment="PATH=/opt/telco-generator/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# GCS bucket for metrics data - UPDATE THIS WITH YOUR BUCKET NAME
Environment="GCS_RETAIL_BUCKET=gs://gcp-sandbox-field-eng-telco-snmp/metrics/"

# Run the retail data generator:
# - Events go to local SFTP directory
# - Metrics go to GCS bucket (JSON format)
ExecStart=/opt/telco-generator/venv/bin/python3 /opt/telco-generator/retail_data_generator.py \
    --events-dir /sftp/retail/events \
    --metrics-gcs-bucket ${GCS_RETAIL_BUCKET} \
    --files-per-minute 500 \
    --log-level INFO

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=retail-generator

# Resource limits
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target